# alloy/rootfs/etc/s6-overlay/s6-rc.d/init-alloy/run
#!/usr/bin/with-contenv bash
set -e

CONFIG_DIR="/etc/alloy"
CONFIG_FILE="${CONFIG_DIR}/config.alloy"
OPTIONS_FILE="/data/options.json"
SSL_DIR="/ssl"

mkdir -p "${CONFIG_DIR}" /data/alloy

# Read options
LOKI_URL=$(jq -r '.loki_url' "${OPTIONS_FILE}")
LOG_LEVEL=$(jq -r '.log_level // "info"' "${OPTIONS_FILE}")
ADDITIONAL_CONFIG=$(jq -r '.additional_config // ""' "${OPTIONS_FILE}")
BASIC_AUTH_USER=$(jq -r '.basic_auth_user // ""' "${OPTIONS_FILE}")
BASIC_AUTH_PASSWORD=$(jq -r '.basic_auth_password // ""' "${OPTIONS_FILE}")
TLS_CA_CERT=$(jq -r '.tls_ca_cert // ""' "${OPTIONS_FILE}")
TLS_CLIENT_CERT=$(jq -r '.tls_client_cert // ""' "${OPTIONS_FILE}")
TLS_CLIENT_KEY=$(jq -r '.tls_client_key // ""' "${OPTIONS_FILE}")
TLS_SERVER_NAME=$(jq -r '.tls_server_name // ""' "${OPTIONS_FILE}")

# Apply extra_hosts entries to /etc/hosts
EXTRA_HOSTS_COUNT=$(jq -r '.extra_hosts // [] | length' "${OPTIONS_FILE}")
if [ "${EXTRA_HOSTS_COUNT}" -gt 0 ]; then
    jq -r '.extra_hosts[]' "${OPTIONS_FILE}" | while IFS=: read -r hostname ip; do
        if [ -n "${hostname}" ] && [ -n "${ip}" ]; then
            echo "${ip} ${hostname}" >> /etc/hosts
        fi
    done
fi

# Detect journal path
JOURNAL_PATH="/var/log/journal"
if [ ! -d "${JOURNAL_PATH}" ] || [ -z "$(ls -A "${JOURNAL_PATH}" 2>/dev/null)" ]; then
    JOURNAL_PATH="/run/log/journal"
fi

# Validate certificate files exist
HAS_TLS=false
if [ -n "${TLS_CA_CERT}" ] || [ -n "${TLS_CLIENT_CERT}" ]; then
    HAS_TLS=true
    for cert_var in TLS_CA_CERT TLS_CLIENT_CERT TLS_CLIENT_KEY; do
        cert_file="${!cert_var}"
        if [ -n "${cert_file}" ] && [ ! -f "${SSL_DIR}/${cert_file}" ]; then
            echo "ERROR: ${cert_var}=${cert_file} not found in ${SSL_DIR}/"
            echo "Place the file in your Home Assistant /ssl/ directory."
            exit 1
        fi
    done
fi

HAS_BASIC_AUTH=false
if [ -n "${BASIC_AUTH_USER}" ]; then
    HAS_BASIC_AUTH=true
fi

# Startup banner
echo "-----------------------------------------------------------"
echo " Grafana Alloy for Home Assistant"
echo " Loki URL: ${LOKI_URL}"
echo " Journal path: ${JOURNAL_PATH}"
echo " Log level: ${LOG_LEVEL}"
if [ "${HAS_TLS}" = true ]; then
    echo " TLS: enabled (CA: ${TLS_CA_CERT:-none}, cert: ${TLS_CLIENT_CERT:-none}, SNI: ${TLS_SERVER_NAME:-auto})"
fi
if [ "${HAS_BASIC_AUTH}" = true ]; then
    echo " Basic auth: enabled (user: ${BASIC_AUTH_USER})"
fi
echo "-----------------------------------------------------------"

# Build loki.write endpoint block
ENDPOINT_BLOCK="    url = \"${LOKI_URL}\""

if [ "${HAS_TLS}" = true ]; then
    TLS_BLOCK="
    tls_config {"
    if [ -n "${TLS_CA_CERT}" ]; then
        TLS_BLOCK="${TLS_BLOCK}
      ca_file   = \"${SSL_DIR}/${TLS_CA_CERT}\""
    fi
    if [ -n "${TLS_CLIENT_CERT}" ]; then
        TLS_BLOCK="${TLS_BLOCK}
      cert_file = \"${SSL_DIR}/${TLS_CLIENT_CERT}\""
    fi
    if [ -n "${TLS_CLIENT_KEY}" ]; then
        TLS_BLOCK="${TLS_BLOCK}
      key_file  = \"${SSL_DIR}/${TLS_CLIENT_KEY}\""
    fi
    if [ -n "${TLS_SERVER_NAME}" ]; then
        TLS_BLOCK="${TLS_BLOCK}
      server_name = \"${TLS_SERVER_NAME}\""
    fi
    TLS_BLOCK="${TLS_BLOCK}
    }"
    ENDPOINT_BLOCK="${ENDPOINT_BLOCK}${TLS_BLOCK}"
fi

if [ "${HAS_BASIC_AUTH}" = true ]; then
    ENDPOINT_BLOCK="${ENDPOINT_BLOCK}

    basic_auth {
      username = \"${BASIC_AUTH_USER}\"
      password = \"${BASIC_AUTH_PASSWORD}\"
    }"
fi

# Generate Alloy config
cat > "${CONFIG_FILE}" <<ALLOYCONFIG
// Grafana Alloy config - auto-generated by HAOS add-on
// Do not edit directly - configure via add-on options

logging {
  level = "${LOG_LEVEL}"
}

// --- Read systemd journal ---
loki.source.journal "journal" {
  path         = "${JOURNAL_PATH}"
  forward_to   = [loki.process.journal.receiver]
  relabel_rules = loki.relabel.journal.rules
  labels       = {
    job = "systemd-journal",
  }
}

// --- Relabel journal fields ---
loki.relabel "journal" {
  forward_to = []

  rule {
    source_labels = ["__journal__systemd_unit"]
    target_label  = "unit"
  }
  rule {
    source_labels = ["__journal__hostname"]
    target_label  = "hostname"
  }
  rule {
    source_labels = ["__journal_syslog_identifier"]
    target_label  = "syslog_identifier"
  }
  rule {
    source_labels = ["__journal__transport"]
    target_label  = "transport"
  }
  rule {
    source_labels = ["__journal_container_name"]
    target_label  = "container_name"
  }
  rule {
    source_labels = ["__journal_priority_keyword"]
    target_label  = "level"
  }
}

// --- Process journal logs ---
loki.process "journal" {
  stage.drop {
    expression = "^\\\\s*$"
  }

  forward_to = [loki.write.loki.receiver]
}

// --- Write to Loki ---
loki.write "loki" {
  endpoint {
${ENDPOINT_BLOCK}
  }
}
ALLOYCONFIG

# Append additional config if provided
if [ -n "${ADDITIONAL_CONFIG}" ] && [ "${ADDITIONAL_CONFIG}" != "null" ]; then
    echo "" >> "${CONFIG_FILE}"
    echo "// --- Additional user config ---" >> "${CONFIG_FILE}"
    echo "${ADDITIONAL_CONFIG}" >> "${CONFIG_FILE}"
fi

echo "Alloy config generated at ${CONFIG_FILE}"
